# scripts/backup.service
[Service]
ExecStart=/home/pineman/services/pinecone-scripts/backup.sh

[Install]
WantedBy=multi-user.target
