# scripts/backup.service
[Service]
ExecStart=/home/pineman/services/scripts/backup.sh

[Install]
WantedBy=multi-user.target
